# src/engram/templates/engram-surface-daemon.service.template
[Unit]
Description=ENGRAM surface daemon (semantic recall)

[Service]
Type=simple
ExecStart={{PYTHON}} {{ENGRAM_HOOKS_DIR}}/engram-surface-daemon.py
Environment=ENGRAM_HOME={{ENGRAM_HOME}}
Environment=ENGRAM_DAEMON_IDLE_TIMEOUT=0
Restart=always
RestartSec=5s
# `append:` log-routing requires systemd >= 240 (2018-10-03).
# Pre-v240 systemd (Ubuntu 18.04 ships 237) would fail to parse this unit.
# Practical floor is fine — engram-alpha's Python 3.10+ floor already
# excludes 18.04 (3.10 not in main repos). Switch to `StandardOutput=journal`
# if a pre-v240 target ever needs support.
StandardOutput=append:{{ENGRAM_HOME}}/surface-daemon.log
StandardError=append:{{ENGRAM_HOME}}/surface-daemon.log

[Install]
WantedBy=default.target
